# etc/systemd/system/fan_auto.service
[Unit]
Description=Run fan_auto

[Service]
Type=oneshot
ExecStart=/usr/bin/sudo /usr/local/bin/fan_auto
